# whatsapp_bot/aremu-reminder.service
[Unit]
Description=Aremu WhatsApp 24-Hour Reminder Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/whatsapp_bot
Environment=PATH=/path/to/whatsapp_bot/venv/bin
Environment=PYTHONPATH=/path/to/whatsapp_bot
ExecStart=/path/to/whatsapp_bot/venv/bin/python reminder_daemon.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aremu-reminder

# Environment variables (set these in production)
Environment=WHATSAPP_ACCESS_TOKEN=your_token_here
Environment=WHATSAPP_PHONE_NUMBER_ID=your_phone_id_here
Environment=DATABASE_URL=your_database_url_here
Environment=OPENAI_API_KEY=your_openai_key_here

[Install]
WantedBy=multi-user.target
